# challenges/Systeme/graduated/crontab
SHELL=/bin/bash

* * * * * /usr/bin/python3 /home/rector/integrator.py >> /home/rector/integrator.log 2>&1
